# ovs-debug/etc/profile.d/vagrant.sh
# Add scripts in vagarnt shared dir to path
export PATH=/vagrant/bin:"${PATH}"
